# recursor/supervisord-recursor.conf
[supervisord]
nodaemon=true

[program:pdns]
command=/usr/sbin/pdns_recursor

